# config/systemd/laravel-pulse.service
[Unit]
Description=Laravel Pulse Server Monitor
After=network.target

[Service]
Type=simple
User=apache
Group=apache
WorkingDirectory=/var/www
ExecStart=/usr/bin/php /var/www/artisan pulse:check
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=laravel-pulse

# Environment
Environment=PATH=/usr/bin:/usr/local/bin
Environment=APP_ENV=production
EnvironmentFile=/var/www/.env

# Resource limits
LimitNOFILE=4096
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
